# root/etc/s6-overlay/s6-rc.d/init-embystat-config/run
#!/usr/bin/with-contenv bash

# make our folders
mkdir -p /config/logs

# permissions
lsiown -R abc:abc \
    /config
