# root/etc/services.d/embystat/run
#!/usr/bin/with-contenv bash

cd /opt/embystat || exit

exec \
	s6-setuidgid abc /opt/embystat/EmbyStat --no-updates \
		--data-dir /config --log-dir /config/logs --config-dir /config
